# docker/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/tmp/supervisord.pid

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true

[program:app]
command=node /opt/app/dist/main.js
directory=/opt/app
autostart=true
autorestart=true
environment=NODE_ENV=production
environment=PORT=3000

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock